# contrib/startup/gentoo/netxmsd
#!/sbin/runscript
# Copyright 1999-2006 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /cvs/nms/contrib/startup/gentoo/netxmsd,v 1.1 2006-01-06 15:41:25 victor Exp $

depend() {
	need net
}

start() {
	ebegin "Starting NetXMS server"
	eval start-stop-daemon --start --quiet --pidfile /var/run/netxmsd.pid --exec /usr/local/bin/netxmsd -- --pid-file /var/run/netxmsd.pid
	eend $?
}

stop() {
	ebegin "Stopping NetXMS server"
	start-stop-daemon --stop --quiet --pidfile /var/run/netxmsd.pid
	eend $?
}
