# packages/linux/debian/netxms-server.postrm
#!/bin/sh

if test "x$1" = "xpurge"; then
	update-rc.d netxmsd remove >/dev/null
	rm -f /etc/netxmsd.conf
	rm -f /etc/netxmsd.conf.backup
	rm -f /etc/netxmsd.conf.sample
	rm -rf /usr/share/netxms
fi
