# packages/linux/debian/6/netxms-agent.postinst
#!/bin/sh

if [ -x /etc/init.d/nxagentd ]; then
	update-rc.d nxagentd defaults 98 2 >/dev/null || exit $?
	if [ -f /etc/nxagentd.conf ]; then
		/etc/init.d/nxagentd start
	else
		echo "****************************************************"
		echo "WARNING"
		echo "Configuration file /etc/nxagentd.conf not found"
		echo "Please create it and then start agent by running"
		echo "/etc/init.d/nxagentd start"
		echo "****************************************************"
	fi
fi
