# contrib/startup/systemd/netxms-server.service.in
[Unit]
Description=NetXMS core server
# StartLimitIntervalSec=0 # since systemd v230

[Service]
Type=simple
ExecStart=~BINDIR~/netxmsd -S
PIDFile=/var/run/netxmsd.pid
LimitNOFILE=65535
LimitCORE=infinity
Restart=on-failure
RestartSec=30
TimeoutSec=900
StartLimitInterval=0 # depricated in v230

[Install]
WantedBy=multi-user.target
Alias=netxmsd.service
