# contrib/startup/redhat/nxhttpd.in
#!/bin/sh
# $Id: nxhttpd.in,v 1.2 2007-09-21 07:37:49 victor Exp $
#
# nxhttpd This script starts and stops the nxhttpd daemon
#
# chkconfig: 2345 90 10
#
# description: nxhttpd is a NetXMS web server process

NETXMS_BINDIR="~BINDIR~"
NETXMS_LIBDIR="~LIBDIR~"

# Source function library.
if [ -f /etc/init.d/functions ]; then
	. /etc/init.d/functions
elif [ -f /etc/rc.d/init.d/functions ]; then
	. /etc/rc.d/init.d/functions
else
	echo "Could not find functions file, your system may be broken."
	exit 1
fi

# Source networking configuration.
. /etc/sysconfig/network

# Check that networking is up.
[ ${NETWORKING} = "no" ] && exit 0

# Default config
NXHTTPD_OPTIONS=""
NXHTTPD_CONFIG="/etc/nxhttpd.conf"
NXHTTPD_PID_FILE="/var/run/nxhttpd.pid"

# Source nxagentd configuration.
if [ -f /etc/sysconfig/nxhttpd ] ; then
        . /etc/sysconfig/nxhttpd
fi

[ -f $NETXMS_BINDIR/nxhttpd ] || exit 0
LD_LIBRARY_PATH=$LD_LIBRARY_PATH:$NETXMS_LIBDIR

# See how we were called.
case "$1" in
  start)
	# Start daemon.
	echo -n "Starting NetXMS Web Server: "
	daemon $NETXMS_BINDIR/nxhttpd -d $NXHTTPD_OPTIONS -c $NXHTTPD_CONFIG -p $NXHTTPD_PID_FILE
	RETVAL=$?
        echo
        [ $RETVAL = 0 ] && touch /var/lock/subsys/nxhttpd
        ;;
  stop)
        # Stop daemons.
        echo -n "Shutting down NetXMS Web Server: "
        kill `cat $NXHTTPD_PID_FILE 2>/dev/null` 2> /dev/null
        RETVAL=$?
        if [ $RETVAL = 0 ] ; then
                success
        else
                failure
        fi
        echo
        [ $RETVAL = 0 ] && rm -f /var/lock/subsys/nxhttpd
        ;;
  restart)
        $0 stop
        $0 start
        ;;
  condrestart)
       [ -e /var/lock/subsys/nxagentd ] && $0 restart
       ;;
  status)
	status nxhttpd
	;;
  *)
	echo "Usage: $0 {start|stop|restart|status|condrestart}"
	exit 1
esac

exit 0
